# docker-images/sds-cron-tasks.prod.crontab
# remember to remove docker image for sds when it was updated, docker will not pull the newest one automatically
*/20 * * * * PATH=$PATH:/snap/bin SDS_SKIP_REENRICH='jira,gerrit,confluence,bugzilla' SDS_SILENT=1 SDS_NCPUS_SCALE=2 SDS_ES_BULKSIZE=100 SDS_TASKS_SKIP_RE='(hyperledger-indy-jira|hyperledger-indy-github|hyperledger-sawtooth-jira|lf-stackstorm-github|lfai-horovod-github|lfai-milvus-github|lfai-onnx-github|lfai-pyro-github|lfn-opnfv-confluence|lfn-onap-confluence|dronecode-github|hyperledger-fabric-jira|lfai-acumos-gerrit|lfedge-akraino-edge-stack-gerrit|lfn-fdio-gerrit|lfn-fdio-jira|lfn-odl-gerrit|lfn-odl-jira|lfn-onap-gerrit|lfn-onap-jira|lfn-opnfv-gerrit|lfn-opnfv-jira|lfn-tungsten-fabric-gerrit|o-ran-shared-gerrit|zephyr-github)' /usr/bin/sds-cron-task.sh sds_main prod 1>> /tmp/sds_main.log 2>>/tmp/sds_main.err
*/10 * * * * PATH=$PATH:/snap/bin SDS_SKIP_REENRICH='jira,gerrit,confluence,bugzilla' SDS_CSV_PREFIX='/root/.perceval/slow_jobs' SDS_SILENT=1 SDS_TASK_TIMEOUT_SECONDS=43200 SDS_ES_BULKSIZE=50 SDS_ONLY_P2O=1 SDS_TASKS_RE='(hyperledger-indy-jira|hyperledger-indy-github|hyperledger-sawtooth-jira|lf-stackstorm-github|lfai-horovod-github|lfai-milvus-github|lfai-onnx-github|lfai-pyro-github|lfn-opnfv-confluence|lfn-onap-confluence|dronecode-github|hyperledger-fabric-jira|lfai-acumos-gerrit|lfedge-akraino-edge-stack-gerrit|lfn-fdio-gerrit|lfn-fdio-jira|lfn-odl-gerrit|lfn-odl-jira|lfn-onap-gerrit|lfn-onap-jira|lfn-opnfv-gerrit|lfn-opnfv-jira|lfn-tungsten-fabric-gerrit|o-ran-shared-gerrit|zephyr-github)' /usr/bin/sds-cron-task.sh sds_slow prod 1>> /tmp/sds_slow.log 2>>/tmp/sds_slow.err
5 * * * * /usr/bin/docker system prune -f 1>>/tmp/prune.log 2>>/tmp/prune.err
