# docker-images/sds.crontab
0 0,6,12,18 * * * /usr/local/bin/sds_cron.sh 1>>/tmp/sds_cron.log 2>>/tmp/sds_cron.err
